# roles/tunnel/templates/matrix-tunnel.service.j2
[Unit]
Description=Matrix federation SSH tunnel (Iran <-> EU)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/ssh -N -o ExitOnForwardFailure=yes \
    -o ServerAliveInterval=15 \
    -o ServerAliveCountMax=3 \
    -o ConnectTimeout=10 \
    -o StrictHostKeyChecking=accept-new \
    -o ProxyCommand="nc -X 5 -x 127.0.0.1:{{ tunnel_socks_port }} %%h %%p" \
    -L {{ tunnel_eu_federation_port }}:127.0.0.1:{{ tunnel_eu_federation_port }} \
    -R {{ matrix_federation_port }}:127.0.0.1:{{ matrix_federation_port }} \
    -i {{ tunnel_ssh_key_path }} \
    root@{{ tunnel_eu_ip }}
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
